# roles/denodo/templates/denodo-init.j2
#! /bin/sh -e
### BEGIN INIT INFO
# Provides: denodo-vql
# Required-Start: $remote_fs $syslog
# Required-Stop: $remote_fs $syslog
# Default-Start: 2 3 4 5
# Default-Stop: 0 1 6
# Short-Description: denodo-vql
# Description: Run Denodo VQL server
# 
### END INIT INFO
# Generated by Ansible for {{ansible_fqdn}}

PATH="/sbin:/bin:/usr/sbin:/usr/bin" #Ne pas toucher

name="denodo-vql"
pidfile="/var/run/$name.pid"
WorkingDirectory=/home/mip/denodo-5.5/
DAEMON="/usr/bin/hugo"
daemon_OPT="server --buildDrafts --watch --baseURL=http://hbps1.intranet.chuv:7011/sp8-documentation/ --bind 0.0.0.0 --port 7011"
daemon_NAME="/home/mip/denodo-5.5/jre/bin/java"
User=mip
Group=mip

. /lib/lsb/init-functions

d_start () {
        log_daemon_msg "Starting system $name Daemon"
	start-stop-daemon --background --pidfile $pidfile --name $daemon_NAME --start --quiet --chuid $User --chdir $WorkingDirectory --exec $DAEMON -- $daemon_OPT
        log_end_msg $?
}

d_stop () {
        log_daemon_msg "Stopping system $name Daemon"
        start-stop-daemon --pidfile $pidfile --name $daemon_NAME --stop --retry 5 --quiet
	log_end_msg $?
}

case "$1" in

        start|stop)
                d_${1}
                ;;

        restart|reload|force-reload)
                        d_stop
                        d_start
                ;;

        force-stop)
               d_stop
                killall -q $daemon_NAME || true
                sleep 2
                killall -q -9 $daemon_NAME || true
                ;;

        status)
                status_of_proc "$daemon_NAME" "$DAEMON" "system-wide $daemon_NAME" && exit 0 || exit $?
                ;;
        *)
                echo "Usage: /etc/init.d/$daemon_NAME {start|stop|force-stop|restart|reload|force-reload|status}"
                exit 1
                ;;
esac
exit 0
